# kubernetes/profile.d/kubernetes.sh
export PATH=/opt/kubernetes/platforms/linux/amd64:$PATH
